# deploy/systemd/miqraa.service.d/override.conf
# Drop-in for `miqraa` systemd unit — install under:
#   /etc/systemd/system/miqraa.service.d/override.conf
# or merge via: sudo systemctl edit miqraa
#
# Raises RLIMIT_MEMLOCK so mediasoup can use io_uring zero-copy buffers without
# falling back; raises file descriptor limit for many WebRTC transports.

[Service]
LimitMEMLOCK=infinity
LimitNOFILE=65536
